# monitor_test.service
[Unit]
Description=Monitor Test Process
After=network.target

[Service]
ExecStart=/usr/local/bin/monitor_test.sh
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target
